# tools/server/hetzner-backup.service
[Unit]
Description=Mirror projects to Hetzner Storage Box

[Service]
Type=oneshot
User=dev
ExecStart=/bin/bash /home/dev/projects/tools/server/rclone-backup.sh
Nice=10
IOSchedulingClass=idle
