# tools/mac/com.hetzner.mounts.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- Auto-mount Hetzner volumes at login.
     Install:  cp tools/mac/com.hetzner.mounts.plist ~/Library/LaunchAgents/
               launchctl load ~/Library/LaunchAgents/com.hetzner.mounts.plist
     Edit the path below to where this repo lives on your Mac. -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.hetzner.mounts</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/milesdirmann/infra/tools/mac/mount-hetzner.sh</string>
    <string>mount</string>
  </array>
  <key>RunAtLoad</key><true/>
  <!-- Self-healing: re-check every 5 min and remount anything that dropped
       (sleep/wake, Wi-Fi change, VPN flap). mount is idempotent when healthy. -->
  <key>StartInterval</key><integer>300</integer>
  <key>StandardErrorPath</key><string>/tmp/hetzner-mounts.err</string>
  <key>EnvironmentVariables</key>
  <dict>
    <!-- /sbin + /usr/sbin required: mount/umount/diskutil live there -->
    <key>PATH</key><string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/sbin:/usr/sbin</string>
  </dict>
</dict>
</plist>
